# deploy/ninesense-guestbook.service
[Unit]
Description=NineSense guestbook API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ninesense
Group=ninesense
WorkingDirectory=/opt/ninesense-guestbook/current
EnvironmentFile=/etc/ninesense/guestbook.env
ExecStart=/opt/ninesense-guestbook/current/venv/bin/uvicorn ninesense_guestbook.app:create_app --factory --host 127.0.0.1 --port 8812 --workers 1 --proxy-headers --forwarded-allow-ips=127.0.0.1 --no-access-log
Restart=on-failure
RestartSec=3
TimeoutStopSec=20
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/ninesense
ReadOnlyPaths=/etc/ninesense/guestbook.env /opt/ninesense-guestbook

[Install]
WantedBy=multi-user.target

